# files/etc/uci-defaults/cee-luci-default-settings
#!/bin/sh

set luci.themes.Argonne=/luci-static/argonne
set luci.main.mediaurlbase=/luci-static/argonne
